# crontab
@hourly /usr/bin/python /home/jake/metadataschema/get_metadata_sheet.py
